# config/jes.crontab
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.XXXXSDQPfd installed on Sun Feb 12 21:35:41 2012)
# (Cron version V5.0 -- $Id: crontab.c,v 1.12 2004/01/23 18:56:42 vixie Exp $)
26 * * * * /home/jes/dumpdb.pl 2>&1 | mail -s cron jes@j-e-s.net
56 3 * * * /home/jes/rsyncdb 2>&1 | mail -s rsyncdb jes@j-e-s.net
16 4 * * * /home/jes/rsyncsite 2>&1 | mail -s rsyncsite jes@j-e-s.net
16 0 1 * * /bin/cat /home/jes/age-products | /usr/bin/psql -U jes voko 2>&1 | mail -s age-products jes@j-e-s.net
