# config/root.crontab
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.XXXXo38Kos installed on Fri Jan  4 13:34:47 2013)
# (Cron version V5.0 -- $Id: crontab.c,v 1.12 2004/01/23 18:56:42 vixie Exp $)
0 0 * * * /root/save-site.pl
0,30 * * * * /etc/exim/make_mail_cdbs

